# packaging/systemd/natmeg-bidsifier.service
[Unit]
Description=NatMEG-BIDSifier GUI (Xvfb + noVNC optional)
After=network.target

[Service]
Type=simple
User=%i
PAMName=login
WorkingDirectory=/opt/NatMEG-BIDSifier
Environment=DISPLAY=:99
ExecStart=/opt/NatMEG-BIDSifier/launch-with-xvfb.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
